# etc/systemd/system/wepn-poweroff.service
[Unit]
Description=WEPN clearnig state at shutdown
DefaultDependencies=no
After=final.target

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/bin/true
ExecStop=/usr/bin/python3 /usr/local/pproxy/system_services/blank_state.py

[Install]
WantedBy=final.target
