# etc/network/if-up.d/wepn
#!/bin/sh
#send heartbeat/correct the LCD status
/usr/bin/python3 /usr/local/pproxy/periodic/send_heartbeat.py
